# src/packages/Louvre/louvre.service
[Unit]
Description=Louvre Display Manager
Conflicts=getty@tty1.service
After=plymouth-quit.service systemd-user-sessions.service

[Service]

UtmpIdentifier=tty1
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

Environment="XDG_RUNTIME_DIR=/tmp"
ExecStart=/bin/louvre.sh
Restart=always
# BusName=org.freedesktop.DisplayManager

[Install]
WantedBy=multi-user.target